# base/systemd-unit-files/ticketweb-authsystem-prunesessions@.service.d/user-install.conf
[Service]



User=
Group=

ExecStart=
ExecStart=%h/.local/libexec/ticketweb-authsystem-prunesessions --devel-mode %f
